# django/app/gunicorn.ini
[program:gunicorn]
command=/usr/local/bin/gunicorn sitio_web.wsgi:mi_aplicacion  --bind 0.0.0.0:8000
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true